# .config/systemd/user/duckdns.service
[Unit]
Description=DuckDNS updater

[Service]
ExecStart=/home/enzime/bin/duckdns
