# .config/systemd/user/duckdns.service
[Unit]
Description=DuckDNS updater

[Service]
ExecStart=/usr/local/bin/duckdns
